# docker/etc/cron.d/bpmn-timer
* * * * * cd /var/www/html/ && php artisan bpmn:timmer >> /dev/null 2>&1
